# web/supervisor/carreraweb.conf
[program:carreraweb]
command=/home/pi/carrera/bin/run_web
autostart=true
autorestart=true

