# web/supervisor/carreraserial.conf
[program:carreraserial]
command=/home/pi/carrera/bin/run_serial
autostart=true
autorestart=true
redirect_stderr=true
